# puppet/vendor/couchdb/files/upstart.conf
# Upstart file at /etc/init/couchdb.conf
# CouchDB
 
start on runlevel [2345]
stop on runlevel [06]
 
pre-start script
    test -d /var/run/couchdb/ || mkdir /var/run/couchdb/
    chown -R couchdb /var/lib/couchdb
    chown -R couchdb /var/log/couchdb
    chown -R couchdb /var/run/couchdb
end script
 
script
    if [ -f /etc/default/couchdb ]; then
        . /etc/default/couchdb
    fi
    exec sudo -u couchdb -E --login /usr/bin/couchdb
end script

post-start script
    sleep 3
end script
 
respawn
respawn limit 10 5
